# mtproxy-bot.service
[Unit]
Description=MTProxy Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mtproxy-bot
EnvironmentFile=/opt/mtproxy-bot/.env
ExecStart=/opt/mtproxy-bot/.venv/bin/python main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
